# host-jez-pc-01/config/systemd/user/cps.service
[Unit]
Description=Calibre-Web

[Service]
Type=simple
ExecStart=/home/jez/bin/books

[Install]
WantedBy=default.target
